# config/init/uwsgi.conf
# simple uWSGI script

description "{{ project_name }} uwsgi instance"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --die-on-term --ini /site/config/uwsgi.ini
